# gunicorn_lingyuanyan.socket
[Unit]
Description=Lingyuan Yan socket
[Socket]
ListenStream= /home/wwwadmin/lingyuanyan/lingyuanyan.sock 
[Install]
WantedBy=sockets.target
